# utils/deadman_switch.service
[Unit]
Description=Wireguard Deadman Switch

[Service]
Type=simple
ExecStart=/usr/bin/deadman_switch.sh

[Install]
WantedBy=default.target
